# launchd/com.ellingrov.amazon-job-radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ellingrov.amazon-job-radar</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>export PATH="/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"; cd /Users/ellingrov/amazon-job-radar-runner; set -a; source /Users/ellingrov/amazon-job-radar-runner/.env.local; set +a; echo "[$(date '+%Y-%m-%d %H:%M:%S')] Starting Amazon job check" &gt;&gt; /Users/ellingrov/amazon-job-radar-runner/output/local-runner/amazon-job-radar.log; npm run check &gt;&gt; /Users/ellingrov/amazon-job-radar-runner/output/local-runner/amazon-job-radar.log 2&gt;&amp;1; echo "[$(date '+%Y-%m-%d %H:%M:%S')] Amazon job check completed" &gt;&gt; /Users/ellingrov/amazon-job-radar-runner/output/local-runner/amazon-job-radar.log</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ellingrov/amazon-job-radar-runner</string>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ellingrov/amazon-job-radar-runner/output/local-runner/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ellingrov/amazon-job-radar-runner/output/local-runner/launchd.stderr.log</string>
  </dict>
</plist>
